# examples/hello/usr/lib/systemd/system/hello-pre.target

[Unit]
Description=Preparation for Shoeshiner
Requires=basic.target
After=basic.target
AllowIsolate=yes
